# varnish_supervisor.conf
[program:varnish]
command=varnishd -F -a 0.0.0.0:$LISTEN_PORT -f $VCL_CONFIG
